# rootfs/etc/services.d/nslcd/run
#!/usr/bin/with-contenv sh
# shellcheck shell=sh

nslcd --version

if [ -n "${DEBUG}" ]; then
    NSLCD_OPTS="${NSLCD_OPTS:-} --debug"
fi

exec nslcd --nofork ${NSLCD_OPTS:-}
