# rootfs/etc/services.d/sshd/run
#!/usr/bin/with-contenv sh
# shellcheck shell=sh

exec /usr/bin/sshd -D -p 2222 ${OPTS:-}
