# rootfs/etc/services.d/sshd/run
#!/usr/bin/with-contenv sh
# shellcheck shell=sh

exec /usr/sbin/sshd -D -p 2222 -h /data/ssh_host_ed25519_key ${OPTS:-}
