# rootfs/etc/services.d/sshd/run
#!/usr/bin/with-contenv sh

exec sshd -p 2222 ${OPTS:-}
